# .config/systemd/user/joplinsync.service
[Unit]
Description=Sync joplin
After=network.target[Service]

[Service]
Type=oneshot
ExecStart=joplin sync

[Install]
WantedBy=default.target
