# .config/systemd/user/joplinsync.timer
[Unit]
Description=Sync joplin

[Timer]
OnBootSec=15min
OnCalendar=*-*-* *:00:00

[Install]
WantedBy=timers.target
