# .config/systemd/user/tasksync.service
[Unit]
Description=Sync taskwarrior 
After=network.target[Service]

[Service]
Type=oneshot
ExecStart=task sync

[Install]
WantedBy=default.target
